# dist/mastodon-web.service
[Unit]
Description=ephemeral.glitch.social-web
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/home/mastodon/ephemeral.glitch.social
Environment="RAILS_ENV=production"
Environment="PORT=3001"
ExecStart=/home/mastodon/.rbenv/shims/bundle exec puma -C config/puma.rb
ExecReload=/bin/kill -SIGUSR1 $MAINPID
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
